# root/etc/s6-overlay/s6-rc.d/init-ransynsrv/run
#!/command/with-contenv sh
# ==============================================================================
# RanSynSrv - Initialization Script
# Integrates: stdout-logs, package-install, nvm, claude-code setup
# ==============================================================================

echo "=========================================="
echo " RanSynSrv - Initializing"
echo "=========================================="

# Runtime Package Installation (universal-package-install mod)
if [ -n "$INSTALL_PACKAGES" ]; then
    # Validate package names (alphanumeric, hyphens, underscores, dots, spaces only)
    if echo "$INSTALL_PACKAGES" | grep -qE '[^a-zA-Z0-9._-]'; then
        echo "[init] ERROR: Invalid characters in INSTALL_PACKAGES (only alphanumeric, dots, hyphens, underscores allowed)"
        exit 1
    fi
    echo "[init] Installing Alpine packages: $INSTALL_PACKAGES"
    apk add --no-cache $INSTALL_PACKAGES || { echo "[init] ERROR: Failed to install packages"; exit 1; }
fi

if [ -n "$INSTALL_PIP_PACKAGES" ]; then
    # Validate package names (alphanumeric, hyphens, underscores, dots, spaces only)
    if echo "$INSTALL_PIP_PACKAGES" | grep -qE '[^a-zA-Z0-9._-]'; then
        echo "[init] ERROR: Invalid characters in INSTALL_PIP_PACKAGES (only alphanumeric, dots, hyphens, underscores allowed)"
        exit 1
    fi
    echo "[init] Installing pip packages: $INSTALL_PIP_PACKAGES"
    pip3 install --no-cache-dir $INSTALL_PIP_PACKAGES || { echo "[init] ERROR: Failed to install pip packages"; exit 1; }
fi

# Create directories
mkdir -p /data/nginx /data/webroot/public_html /data/webroot/goaccess /data/databases
mkdir -p /data/log/nginx /data/log/php /data/ssh /data/scripts /data/crontabs
mkdir -p /data/claude/.claude /data/commandhistory /workspace

# Copy defaults
if [ ! -f /data/nginx/nginx.conf ]; then
    echo "[init] Installing default nginx.conf"
    cp /defaults/nginx/nginx.conf /data/nginx/
fi

if [ ! -f /data/webroot/public_html/index.php ] && [ ! -f /data/webroot/public_html/index.html ]; then
    echo "[init] Installing default index.php"
    cp /defaults/webroot/public_html/index.php /data/webroot/public_html/
fi

[ ! -f /data/webroot/goaccess/index.html ] && cp /defaults/webroot/goaccess/index.html /data/webroot/goaccess/

# Link nginx config
ln -sf /data/nginx/nginx.conf /etc/nginx/nginx.conf

# Create log files
touch /data/log/nginx/access.log /data/log/nginx/error.log /data/log/php/error.log

# Stdout Logging (universal-stdout-logs mod)
if [ "$DOCKER_LOGS" = "true" ]; then
    echo "[init] Enabling stdout logging (DOCKER_LOGS=true)"
    rm -f /data/log/nginx/access.log /data/log/nginx/error.log /data/log/php/error.log
    ln -sf /proc/1/fd/1 /data/log/nginx/access.log
    ln -sf /proc/1/fd/2 /data/log/nginx/error.log
    ln -sf /proc/1/fd/2 /data/log/php/error.log
else
    echo "[init] Using file-based logging"
fi

# Set permissions (all persistent data is under /data)
chown -R abc:abc /data /workspace
chmod 700 /data/ssh

echo "=========================================="
echo " RanSynSrv - Ready"
echo "=========================================="
