# root/etc/s6-overlay/s6-rc.d/init-ransynsrv/run
#!/command/with-contenv sh
# ==============================================================================
# RanSynSrv - Initialization (s6-rc oneshot)
# Runs before svc-nginx / svc-php-fpm / svc-goaccess / svc-ttyd via
# dependencies.d/init-ransynsrv on each of those services.
# ==============================================================================

echo "=========================================="
echo " RanSynSrv - Initializing"
echo "=========================================="

SENTINEL=/data/.ransynsrv-init-done

# --- Runtime package install (hash-cached across boots) ----------------------
# apk add without caching would run on every container start (5-15s latency
# when a mirror is slow). Cache a hash of INSTALL_PACKAGES; skip if unchanged.
INSTALL_CACHE_DIR=/data/.ransynsrv-install-cache
mkdir -p "$INSTALL_CACHE_DIR"

if [ -n "$INSTALL_PACKAGES" ]; then
    if echo "$INSTALL_PACKAGES" | grep -qE '[^a-zA-Z0-9._ -]'; then
        echo "[init] ERROR: invalid characters in INSTALL_PACKAGES"
        exit 1
    fi
    HASH=$(printf '%s' "$INSTALL_PACKAGES" | sha256sum | cut -d' ' -f1)
    if [ "$(cat "$INSTALL_CACHE_DIR/apk.hash" 2>/dev/null)" = "$HASH" ]; then
        echo "[init] Alpine packages already installed (hash match): $INSTALL_PACKAGES"
    else
        echo "[init] Installing Alpine packages: $INSTALL_PACKAGES"
        apk add --no-cache $INSTALL_PACKAGES || { echo "[init] ERROR: apk add failed"; exit 1; }
        echo "$HASH" > "$INSTALL_CACHE_DIR/apk.hash"
    fi
fi

if [ -n "$INSTALL_PIP_PACKAGES" ]; then
    if echo "$INSTALL_PIP_PACKAGES" | grep -qE '[^a-zA-Z0-9._ =-]'; then
        echo "[init] ERROR: invalid characters in INSTALL_PIP_PACKAGES"
        exit 1
    fi
    HASH=$(printf '%s' "$INSTALL_PIP_PACKAGES" | sha256sum | cut -d' ' -f1)
    if [ "$(cat "$INSTALL_CACHE_DIR/pip.hash" 2>/dev/null)" = "$HASH" ]; then
        echo "[init] pip packages already installed (hash match): $INSTALL_PIP_PACKAGES"
    else
        echo "[init] Installing pip packages: $INSTALL_PIP_PACKAGES"
        pip3 install --no-cache-dir $INSTALL_PIP_PACKAGES || { echo "[init] ERROR: pip install failed"; exit 1; }
        echo "$HASH" > "$INSTALL_CACHE_DIR/pip.hash"
    fi
fi

# --- Directory scaffolding ---------------------------------------------------
mkdir -p \
    /data/nginx /data/webroot/public_html /data/webroot/goaccess /data/databases \
    /data/log/nginx /data/log/php /data/ssh /data/scripts \
    /data/claude/.claude /data/commandhistory /workspace

# --- First-run defaults ------------------------------------------------------
if [ ! -f /data/nginx/nginx.conf ]; then
    echo "[init] Installing default nginx.conf"
    cp /defaults/nginx/nginx.conf /data/nginx/
fi

if [ ! -f /data/webroot/public_html/index.php ] && [ ! -f /data/webroot/public_html/index.html ]; then
    echo "[init] Installing default index.php"
    cp /defaults/webroot/public_html/index.php /data/webroot/public_html/
fi

[ ! -f /data/webroot/goaccess/index.html ] && cp /defaults/webroot/goaccess/index.html /data/webroot/goaccess/

if [ ! -f /data/CLAUDE.md ]; then
    echo "[init] Installing default CLAUDE.md for Claude Code"
    cp /defaults/CLAUDE.md /data/
fi

# --- PHP timeout template (keeps nginx's fastcgi_read_timeout in sync with
#     PHP_MAX_EXECUTION_TIME so a long-running script can't silently 504) ------
cat > /data/nginx/php-timeout.conf <<TIMEOUT
# Auto-generated by init-ransynsrv from \$PHP_MAX_EXECUTION_TIME. Do not edit;
# include'd by nginx.conf inside the location ~ \.php\$ block.
fastcgi_read_timeout ${PHP_MAX_EXECUTION_TIME:-300};
fastcgi_send_timeout ${PHP_MAX_EXECUTION_TIME:-300};
TIMEOUT

# --- Marker-bounded auth-block writer ---------------------------------------
# Shared helper — replaces the block between <marker>_BEGIN and <marker>_END
# with the given body, atomically (write-then-rename). Warns when the target
# config has no markers (e.g. user customised away our scaffold).
write_marker_block() {
    # $1 = path, $2 = marker-basename (e.g. "GOACCESS_AUTH"), $3 = block body
    python3 - "$1" "$2" "$3" <<'PY'
import os, sys, re
path, marker, body = sys.argv[1], sys.argv[2], sys.argv[3]
with open(path, 'r') as f:
    content = f.read()
begin = f"# {marker}_BEGIN"
end   = f"# {marker}_END"
pattern = re.compile(re.escape(begin) + r'(.*?)' + re.escape(end), re.DOTALL)
if not pattern.search(content):
    print(f"[init] WARNING: {path} has no {marker}_BEGIN / {marker}_END markers; "
          f"skipping auth-block write. Customized nginx.conf? "
          f"Restore the markers to re-enable auth toggling.", file=sys.stderr)
    sys.exit(0)
def _replace(m):
    return begin + "\n" + body + "\n            " + end
new_content = pattern.sub(_replace, content)
if new_content == content:
    sys.exit(0)
tmp = path + ".tmp"
with open(tmp, 'w') as f:
    f.write(new_content)
os.replace(tmp, path)
PY
}

# --- GoAccess auth (nginx-layer Basic auth, credentials hashed at rest) -----
if [ "$GOACCESS_AUTH_ENABLED" = "true" ]; then
    if [ -z "$GOACCESS_PASSWORD" ]; then
        echo "[init] ERROR: GOACCESS_AUTH_ENABLED=true but GOACCESS_PASSWORD is empty"
        exit 1
    fi
    GOACCESS_USERNAME=${GOACCESS_USERNAME:-admin}
    echo "[init] Enabling GoAccess auth for user: ${GOACCESS_USERNAME}"
    # -6 = SHA-512-crypt (~3 orders of magnitude stronger offline-crack
    # resistance than the older -apr1 / MD5-crypt). nginx-http-auth-basic
    # recognizes this format natively.
    GOACCESS_HASH=$(openssl passwd -6 "$GOACCESS_PASSWORD")
    echo "${GOACCESS_USERNAME}:${GOACCESS_HASH}" > /data/nginx/.goaccess-htpasswd
    chmod 600 /data/nginx/.goaccess-htpasswd
    write_marker_block /data/nginx/nginx.conf GOACCESS_AUTH \
'            auth_basic "GoAccess Analytics Dashboard";
            auth_basic_user_file /data/nginx/.goaccess-htpasswd;'
    # Legacy .htpasswd path (pre-multi-auth refactor) — remove if lingering
    [ -f /data/nginx/.htpasswd ] && rm -f /data/nginx/.htpasswd
else
    echo "[init] GoAccess auth disabled"
    write_marker_block /data/nginx/nginx.conf GOACCESS_AUTH '            auth_basic off;'
    [ -f /data/nginx/.goaccess-htpasswd ] && rm -f /data/nginx/.goaccess-htpasswd
    [ -f /data/nginx/.htpasswd ] && rm -f /data/nginx/.htpasswd
fi

# --- ttyd auth (nginx-layer Basic auth; ttyd itself runs without -c) --------
# Moving the credential off ttyd's argv avoids exposing the plaintext password
# in `ps aux` / /proc/<pid>/cmdline to any process running as abc.
if [ "$TTYD_ENABLED" = "true" ] && [ -n "$TTYD_USERNAME" ] && [ -n "$TTYD_PASSWORD" ]; then
    echo "[init] Enabling ttyd nginx-layer auth for user: ${TTYD_USERNAME}"
    TTYD_HASH=$(openssl passwd -6 "$TTYD_PASSWORD")
    echo "${TTYD_USERNAME}:${TTYD_HASH}" > /data/nginx/.ttyd-htpasswd
    chmod 600 /data/nginx/.ttyd-htpasswd
    write_marker_block /data/nginx/nginx.conf TTYD_AUTH \
'            auth_basic "RanSynSrv Terminal";
            auth_basic_user_file /data/nginx/.ttyd-htpasswd;'
else
    if [ "$TTYD_ENABLED" = "true" ]; then
        echo "[init] WARNING: TTYD_ENABLED=true but TTYD_USERNAME and/or TTYD_PASSWORD"
        echo "[init] WARNING: is empty. ttyd is starting WITHOUT authentication."
        echo "[init] WARNING: Set both env vars, or set TTYD_ENABLED=false to disable."
    fi
    write_marker_block /data/nginx/nginx.conf TTYD_AUTH '            auth_basic off;'
    [ -f /data/nginx/.ttyd-htpasswd ] && rm -f /data/nginx/.ttyd-htpasswd
fi

# --- nginx config symlink (user-override path) -------------------------------
# `ln -sfn` replaces an existing symlink atomically on Linux (rename-over
# semantics). The prior `rm -f && ln -s` sequence had a TOCTOU window where a
# crashed init could leave /etc/nginx/nginx.conf missing → nginx crash-loops
# on next boot with no config.
ln -sfn /data/nginx/nginx.conf /etc/nginx/nginx.conf

# --- Log files + DOCKER_LOGS stdout mode -------------------------------------
# Create with mode 0640 (owner+group read, no world) — access logs contain
# real-IP after real_ip_header processing and shouldn't be world-readable.
for f in /data/log/nginx/access.log /data/log/nginx/error.log /data/log/php/error.log; do
    [ -e "$f" ] || install -m 0640 /dev/null "$f"
done

if [ "$DOCKER_LOGS" = "true" ]; then
    echo "[init] Enabling stdout logging (DOCKER_LOGS=true)"
    rm -f /data/log/nginx/access.log /data/log/nginx/error.log /data/log/php/error.log
    ln -sf /proc/1/fd/1 /data/log/nginx/access.log
    ln -sf /proc/1/fd/2 /data/log/nginx/error.log
    ln -sf /proc/1/fd/2 /data/log/php/error.log
else
    echo "[init] Using file-based logging"
fi

# --- .cache for gitstatus/p10k in ttyd zsh sessions --------------------------
mkdir -p /home/abc/.cache
chown -R abc:abc /home/abc/.cache

# --- Ownership + permissions -------------------------------------------------
# First-boot: do a full recursive chown (may touch lots of files if the user
# pre-populated /data). Subsequent boots: only re-chown the small subtrees we
# know the services need to own, avoiding an O(n) stall on large data volumes.
if [ ! -f "$SENTINEL" ]; then
    echo "[init] First boot: recursive chown /data /workspace (one-time)"
    chown -R abc:abc /data /workspace
    touch "$SENTINEL"
else
    chown abc:abc /data /workspace 2>/dev/null || true
    # Targeted: init-managed paths that the services write to on each boot.
    chown -R abc:abc /data/nginx /data/log /data/webroot/goaccess /data/claude /home/abc/.cache 2>/dev/null || true
fi

# SSH keys — dir 0700 and any file inside 0600 (inherited file modes from the
# host bind-mount can be 0644, which ssh-client refuses to use).
chmod 700 /data/ssh
find /data/ssh -mindepth 1 -type f -exec chmod 600 {} + 2>/dev/null || true

echo "=========================================="
echo " RanSynSrv - Ready"
echo "=========================================="
