# root/etc/s6-overlay/s6-rc.d/svc-php-fpm/run
#!/command/with-contenv sh
exec 2>&1

# Regenerate /etc/php84/conf.d/99-ransynsrv.ini from current env on each boot
# so PHP_MEMORY_LIMIT, PHP_MAX_UPLOAD, PHP_MAX_POST, PHP_MAX_EXECUTION_TIME,
# and TZ are actually runtime-configurable (they used to be baked at image
# build time, silently ignoring any .env / compose override).
cat > /etc/php84/conf.d/99-ransynsrv.ini <<INI
date.timezone=${TZ:-UTC}
memory_limit=${PHP_MEMORY_LIMIT:-256M}
upload_max_filesize=${PHP_MAX_UPLOAD:-50M}
post_max_size=${PHP_MAX_POST:-50M}
max_execution_time=${PHP_MAX_EXECUTION_TIME:-300}
max_input_time=${PHP_MAX_EXECUTION_TIME:-300}
expose_php=Off
INI

# Pool config — clear_env=no forwards container env (ANTHROPIC_API_KEY, DB URLs,
# app-defined vars) to PHP so getenv()/$_ENV work in-app. The container is
# single-tenant; if you ever host untrusted PHP here, switch to an explicit
# env[...] allowlist instead.
cat > /etc/php84/php-fpm.d/www.conf <<'POOL'
[www]
user = abc
group = abc
listen = /run/php/php-fpm.sock
listen.owner = abc
listen.group = nginx
listen.mode = 0660
pm = dynamic
pm.max_children = 10
pm.start_servers = 2
pm.min_spare_servers = 1
pm.max_spare_servers = 3
pm.max_requests = 500
php_admin_value[error_log] = /data/log/php/error.log
php_admin_flag[log_errors] = on
clear_env = no
POOL

echo "Starting php-fpm..."
exec php-fpm -F
