# root/etc/s6-overlay/s6-rc.d/svc-ttyd/run
#!/command/with-contenv sh
exec 2>&1

if [ "$TTYD_ENABLED" != "true" ]; then
    echo "ttyd disabled (set TTYD_ENABLED=true to enable)"
    # See svc-goaccess/run for why this is s6-pause, not sleep infinity.
    exec /command/s6-pause
fi

echo "Starting ttyd web terminal on port 7681..."

# Environment for the abc login shell
export HOME=/home/abc
# Gitstatus/p10k need /home/abc/.cache writable; init also does this, but
# ensuring it here makes the service self-contained.
mkdir -p /home/abc/.cache
chown -R abc:abc /home/abc/.cache

cd /data

# Authentication is handled at the nginx layer (see init-ransynsrv's TTYD_AUTH
# marker-block edit + /data/nginx/.ttyd-htpasswd). ttyd binds to lo only
# (`-i lo`), so no external access bypasses the nginx gate, and the credential
# is never in this process's argv / /proc/<pid>/cmdline.
exec s6-setuidgid abc ttyd \
    -p 7681 \
    -i lo \
    -W \
    -t fontSize=14 \
    -t 'theme={"background":"#1e1e2e","foreground":"#cdd6f4"}' \
    /bin/zsh -l
